# sync-ports.service
[Unit]
Description=Sync port scan results to main VPS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root

# ============================================================
# Secrets are stored in /etc/sync-ports.conf (chmod 600)
# — NOT exposed in this file or in ps output.
# ============================================================
EnvironmentFile=/etc/sync-ports.conf
ExecStart=/root/sync-ports.sh

# Auto-restart if it crashes (wait 30s between retries)
Restart=always
RestartSec=30

# Log to journald
StandardOutput=journal
StandardError=journal

# Sanitize logs: hide env var values from journal
# (prevents password leaks in journalctl output)
LogNamespace=sync-ports

[Install]
WantedBy=multi-user.target
